# build/root/etc/services.d/nodejs-web/run
#!/usr/bin/with-contenv bash

npm start --prefix /app/