# root/etc/services.d/wireguard/finish
#!/bin/bash

echo "[$(date -Iseconds)] Connection summary:"
wg show wg0
wg-quick down wg0
rm -f /etc/wireguard/wg0.conf