# root/etc/services.d/wireguard/run
#!/bin/bash

wg-quick up wg0
wg show wg0

sleep infinity &
wait
