# root/etc/services.d/coredns/run
#!/usr/bin/with-contenv bash

cd /config/coredns

exec \
    /app/coredns -dns.port=53