# rootfs/etc/services.d/wireguard/run
#!/bin/bash

_term() {
  wg-quick down wg0
}

trap _term SIGTERM SIGINT EXIT

wg-quick up wg0

sleep infinity &

wait
